# recipes/db/redis-server/redis-supervisord.conf
# supervisord configuration for redis

[program:redis]
command = %PREFIX%/bin/redis-server %PREFIX%/etc/redis.conf
autostart = false
